# system_files/deck/shared/etc/systemd/system/ac.target
[Unit]
Description=On AC power
DefaultDependencies=no
StopWhenUnneeded=yes
