# system_files/desktop/shared/usr/lib/systemd/system/NetworkManager.service.d/iwd.conf
[Unit]
# Ensure iwd is started before NetworkManager when using iwd as WiFi backend
Wants=iwd.service
After=iwd.service
