# system_files/desktop/shared/etc/profile.d/askpass.sh
SUDO_ASKPASS='/usr/bin/ksshaskpass'
export SUDO_ASKPASS
