# system_files/desktop/shared/usr/lib/udev/rules.d/99-cec-bluetooth.rules
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/usr/bin/cec-control onboot"
