# spec_files/unl0kr/unl0kr-ask-password.path
#  SPDX-License-Identifier: LGPL-2.1-or-later

[Unit]
Description=Forward Password Requests to unl0kr Directory Watch
DefaultDependencies=no
After=plymouth-start.service
Before=paths.target cryptsetup.target
Conflicts=emergency.service
Before=emergency.service
Conflicts=shutdown.target
Before=shutdown.target

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes
