# deck/etc/systemd/user/distrobox-upgrade-automatic.service
[Unit]
Description=distrobox-upgrade Automatic Update
 
[Service]
Type=simple
ExecStart=/usr/bin/distrobox-upgrade --all
StandardOutput=null
