# system_files/deck/shared/usr/lib/systemd/system/battery.target
[Unit]
Description=On battery power
DefaultDependencies=no
StopWhenUnneeded=yes
