# system_files/deck/shared/usr/lib/systemd/system/bazzite-grub-boot-success.timer
[Unit]
Description=Mark update as successful after the system has been alive for 3 minutes.

[Timer]
OnStartupSec=3min

[Install]
WantedBy=timers.target
