# system_files/deck/shared/usr/lib/systemd/system/bazzite-autologin.service
[Unit]
Description=Enables Bazzite autologin
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bazzite-autologin

[Install]
WantedBy=multi-user.target
