# system_files/deck/shared/usr/lib/systemd/system/gamescope-autologin.service
[Unit]
Description=Enables autologin to gamescope-session
Before=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/bin/gamescope-autologin

[Install]
WantedBy=multi-user.target
