# system_files/deck/shared/usr/lib/systemd/system/return-to-gamemode.service
[Unit]
Description=Return to Gaming Mode

[Service]
Type=oneshot
ExecStart=/usr/bin/return-to-gamemode
