# system_files/deck/shared/usr/lib/systemd/system/ublue-update-force.service
[Unit]
Description=Force Universal Blue Update Oneshot Service

[Service]
Type=oneshot
ExecStart=/usr/bin/ublue-update --force
RemainAfterExit=no
