# system_files/desktop/shared/usr/lib/systemd/system/bluetooth.service.d/prevent_hanging.conf
# Prevent freezing device with bluetooth controllers
[Unit]
Conflicts=sleep.target
Before=sleep.target