# system_files/deck/shared/usr/lib/systemd/user/bazzite-grub-boot-success.timer
[Unit]
Description=Mark update as successful after the user session has run 3 minutes.
ConditionUser=!@system
ConditionVirtualization=!container

[Timer]
OnActiveSec=3min

[Install]
WantedBy=timers.target
