# system_files/desktop/shared/etc/systemd/system/beesd@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/timeout 30m /usr/bin/beesd --no-timestamps %i
Restart=no
SuccessExitStatus=124
OOMScoreAdjust=1000
OOMPolicy=stop
