# dx/system_files/usr/lib/systemd/system/bazzite-dx-groups.service
[Unit]
Description=Add wheel members to docker, and incus-admin groups

[Service]
Type=oneshot
ExecStart=/usr/libexec/bazzite-dx-groups

[Install]
WantedBy=default.target