# system_files/desktop/kinoite/usr/lib/systemd/system/plasmalogin.service.d/override.conf
# DELETEME: Whenever this fix makes it to fedora.
# See https://invent.kde.org/plasma/plasma-login-manager/-/commit/a8c752fe28c36fa7d3b95461a7b24ce2d9de7f69
[Unit]
ExecStartPre=-udevadm settle --timeout=10
