# system_files/desktop/etc/profile.d/neofetch.sh
alias neofetch='neofetch --source /usr/share/ublue-os/bazzite/logo.txt --config /usr/share/ublue-os/bazzite/neofetch.conf'
