# system_files/desktop/kinoite/usr/lib/systemd/system/sddm.service.d/override.conf
[Unit]
Wants=systemd-udev-settle.service
After=systemd-udev-settle.service
