# system_files/desktop/kinoite/usr/lib/systemd/system/sddm.service.d/override.conf
[Unit]
Wants=plymouth-quit-wait.service
After=plymouth-quit-wait.service
