# system_files/desktop/shared/usr/lib/systemd/system/mcelog.service.d/override.conf
# Override configuration for mcelog.service so that it does not run in virtualized environments.
[Unit]
ConditionVirtualization=no