# system_files/desktop/shared/etc/xdg/autostart/bazzite-announcement.desktop
[Desktop Entry]
Type=Application
Name=Bazzite Announcements
Exec=/usr/libexec/bazzite-announcement
Terminal=false
StartupNotify=false
X-GNOME-Autostart-enabled=true
