# installer/system_files/shared/etc/skel/.config/autostart/on_gui_login.desktop
[Desktop Entry]
Exec=/usr/bin/on_gui_login.sh
Icon=application-x-shellscript
Type=Application
