# system_files/desktop/shared/usr/etc/xdg/autostart/joystickwake.desktop
[Desktop Entry]
Type=Application
Name=Joystick Wake
Comment=Joystick-aware screen waker
Icon=input-gaming
NoDisplay=true
Exec=/usr/bin/bazzite-joystickwake
Terminal=false
