# system_files/desktop/shared/usr/lib/systemd/system/ublue-update.service.d/override.conf
[Service]
Environment="TOPGRADE_SKIP_BRKC_NOTIFY=true"
