# system_files/desktop/kinoite/etc/systemd/user/com.system76.Scheduler.dbusproxy.service
[Unit]
Description=Forward com.system76.Scheduler session DBus messages to the system bus

[Service]
ExecStart=/usr/bin/system76-scheduler-dbus-proxy

[Install]
WantedBy=default.target
