# system_files/desktop/silverblue/etc/systemd/user/sunshine.service.d/override.conf
[Install]
WantedBy=gnome-session.target
