# system_files/deck/shared/usr/lib/systemd/system/hhd@.service.d/override.conf
[Service]
Environment="HHD_QAM_KEYBOARD=1"
