# system_files/desktop/shared/usr/lib/systemd/user/bazzite-auto-pause.service
[Unit]
Description=Bazzite Auto-Pause Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/libexec/bazzite-auto-pause/auto-pause-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
