# system_files/desktop/shared/usr/lib/systemd/user/bazzite-dynamic-fixes.service
[Unit]
Description=Applies conditional user level workarounds
ConditionUser=!@system

[Service]
Type=simple
ExecStart=/usr/libexec/bazzite-dynamic-fixes

[Install]
WantedBy=default.target
