# system_files/desktop/shared/usr/lib/systemd/user/bazzite-user-setup.service
[Unit]
Description=Configure Bazzite for current user

[Service]
Type=simple
ExecStart=/usr/libexec/bazzite-user-setup

[Install]
WantedBy=default.target
