# system_files/desktop/shared/usr/lib/systemd/system/bazzite-flatpak-manager.service
[Unit]
Description=Manage system flatpaks
After=bazzite-hardware-setup.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/bazzite-flatpak-manager
Restart=on-failure
RestartSec=30
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
