# system_files/desktop/shared/usr/lib/systemd/system/bazzite-disable-snapper-workaround.service
# FIXME: Remove me whenever snapper doesnt enable automatic snapshots on config creation

[Unit]
Description=Disable snapper for skewed configs
ConditionPathIsSymbolicLink=/etc/systemd/system/timers.target.wants/snapper-timeline.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c "systemctl disable --now snapper-timeline.timer; systemctl enable --now snapper-cleanup.timer; snapper set-config TIMELINE_LIMIT_YEARLY=0; snapper set-config TIMELINE_LIMIT_MONTHLY=0; systemctl disable %n"

[Install]
WantedBy=multi-user.target