# system_files/desktop/shared/usr/lib/systemd/system/bazzite-snapper-setup.service
[Unit]
Description=Create a default snapper config if snapper is unconfigured
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/bazzite-snapper-setup

[Install]
WantedBy=multi-user.target
