# system_files/desktop/shared/usr/lib/systemd/system/beesd@.timer
[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target