# spec_files/waydroid/dev-binderfs.mount
[Unit]
Description=Mount binderfs partition

[Mount]
What=binder
Where=/dev/binderfs
Type=binder
