# system_files/deck/etc/profile.d/steamos-btrfs.sh
if [[ -f /etc/default/steamos-btrfs ]]; then
	source /etc/default/steamos-btrfs
fi
