# system_files/desktop/shared/usr/etc/profile.d/brew.sh
#!/usr/bin/bash
[[ -d /home/linuxbrew/.linuxbrew && $- == *i* ]] && eval "$(/home/linuxbrew/.linuxbrew/bin/brew shellenv)"
