# system_files/desktop/shared/usr/etc/profile.d/nix-app-icons.sh
XDG_DATA_DIRS="$HOME/.nix-profile/share:/nix/var/nix/profiles/default/share:$XDG_DATA_DIRS"